# ops/macbook/com.lifeos.backup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.lifeos.backup</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/you/life-os/ops/macbook/backup-to-mini.sh</string>
  </array>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key><integer>3</integer>
    <key>Minute</key><integer>0</integer>
  </dict>
  <!-- launchd runs a MISSED 3am job when the Mac next wakes, and never stacks
       duplicates. Combined with the gate (defers on OBS / active / battery),
       a backup can never start while you're recording. -->
  <key>ProcessType</key><string>Background</string>
  <key>LowPriorityIO</key><true/>
  <key>Nice</key><integer>10</integer>
  <key>StandardOutPath</key><string>/Users/you/Library/Logs/lifeos-backup.log</string>
  <key>StandardErrorPath</key><string>/Users/you/Library/Logs/lifeos-backup.log</string>
</dict>
</plist>
