# ops/macbook/com.lifeos.ingest.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.lifeos.ingest</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/you/life-os/ops/macbook/ingest-watch.sh</string>
  </array>
  <key>StartInterval</key><integer>120</integer>
  <key>RunAtLoad</key><true/>
  <!-- politeness: E-cores + throttled disk I/O; interactive work always wins -->
  <key>ProcessType</key><string>Background</string>
  <key>LowPriorityIO</key><true/>
  <key>Nice</key><integer>10</integer>
  <key>StandardOutPath</key><string>/Users/you/Library/Logs/lifeos-ingest.log</string>
  <key>StandardErrorPath</key><string>/Users/you/Library/Logs/lifeos-ingest.log</string>
</dict>
</plist>
