# ops/mini/com.lifeos.nasbackup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.lifeos.nasbackup</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/you/life-os/ops/mini/backup-to-nas.sh</string>
  </array>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key><integer>3</integer>
    <key>Minute</key><integer>30</integer>
  </dict>
  <key>StandardOutPath</key><string>/Users/you/Library/Logs/lifeos-nasbackup.log</string>
  <key>StandardErrorPath</key><string>/Users/you/Library/Logs/lifeos-nasbackup.log</string>
</dict>
</plist>
